# deploy/systemd/campaign-lottery-platform.service
[Unit]
Description=Campaign Lottery Platform API
After=network.target mysql.service redis.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/campaign-lottery-platform/backend
EnvironmentFile=/home/ubuntu/campaign-lottery-platform/backend/.env
ExecStart=/home/ubuntu/campaign-lottery-platform/backend/campaign-lottery-api
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
